# health_check.service
[Unit]
Description=Health check

[Service]
ExecStart=/path/to/health.sh
StandardOutput=journal
StandardError=journal

